# rootfs_files/etc/init.d/installplugins-mos
#!/bin/sh
# MOS custom

### BEGIN INIT INFO
# Provides:          installplugins-mos
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     S
# Default-Stop:
# Short-Description: Install MOS Plugins
# Description:       Run MOS Plugin installation
### END INIT INFO

PATH="/sbin:/bin:/usr/sbin:/usr/bin"
DESC="MOS Plugins"
NAME="installplugins-mos"
RUN_DIR="/run"

. /lib/lsb/init-functions

# Set apt options
APT_OPTS="-o Acquire::Retries=1 -o Acquire::http::Timeout=10 -o Acquire::https::Timeout=10 -o Acquire::ftp::Timeout=10"

# Check if there a plugins with additional packages, if yes execute update
PLG_PKGS=$(cat /boot/optional/plugins/*/packages/packages)
if [ ! -z "$PLG_PKGS" ] ; then
  logger -t "MOS Plugins" "Updating aptitude index"
  timeout 20 env DEBIAN_FRONTEND=noninteractive apt-get.orig $APT_OPTS update >/dev/null 2>&1
fi

# Keep only the newest .deb per package
cleanup_old_debs() {
  local deb base
  for deb in "$1"/*.deb ; do
    [ -e "$deb" ] || continue
    base="${deb##*/}"
    base="${base%%_*}"
    ls -1 "$1/${base}_"*.deb 2>/dev/null | sort -V | head -n -1 | while read -r old ; do
      rm -f "$old"
    done
  done
}

install_plugins() {
  # Get all plugin directories
  for plugin in /boot/optional/plugins/*; do
    # Set installation done and get version paths
    INSTALLATION_DONE=false
    VERSION_PATH=$(ls -d $plugin/*/ 2>/dev/null | grep -E '/[0-9]+\.[0-9]+\.[0-9]+/?$' | sort -V | tail -1)
    if [ ! -z "$VERSION_PATH" ] ; then
      PACKAGE_PATH=$(ls "$VERSION_PATH"*.deb 2>/dev/null)
      # If there is not .deb package in Version folder it's proabably not a plugin
      if [ ! -z "$PACKAGE_PATH" ] ; then
        DISPLAY_NAME=$(jq -r '.name' $plugin/template.json)
        logger -t "MOS Plugins: $DISPLAY_NAME" "Installing: ${PACKAGE_PATH##*/}" 
        DEBIAN_FRONTEND=noninteractive dpkg -i $PACKAGE_PATH >/dev/null 2>&1
        INSTALLATION_DONE=true
      fi
      # Download/Update additional packages when packages file is found
      if [ -f "$plugin/packages/packages" ] ; then
        logger -t "MOS Plugins: $DISPLAY_NAME" "Updating and installing additinal packages" 
        PKGS=$(grep -vE '^[[:space:]]*(#|$)' "$plugin/packages/packages")
        timeout 60 env DEBIAN_FRONTEND=noninteractive apt-get.orig $APT_OPTS install -y --download-only --no-install-recommends \
          -o Dir::Cache::archives="$plugin/packages/" $PKGS >/dev/null 2>&1
        cleanup_old_debs "$plugin/packages"
        timeout 30 env DEBIAN_FRONTEND=noninteractive dpkg -i "$plugin"/packages/*.deb >/dev/null 2>&1
        timeout 15 env DEBIAN_FRONTEND=noninteractive dpkg --configure -a >/dev/null 2>&1
      fi
      # Run functions file if found with timeout to avoid boot blocking
      if [ -f "$VERSION_PATH/functions" ] ; then
        logger -t "MOS Plugins: $DISPLAY_NAME" "Executing start script"
        timeout -s KILL 300 bash -c '. "$1"; type mos_start >/dev/null 2>&1 && mos_start' bash "$VERSION_PATH/functions" 2>&1 | logger -t "MOS Plugins: $DISPLAY_NAME"
        case "${PIPESTATUS[0]}" in
          0)       logger -t "MOS Plugins: $DISPLAY_NAME" "Start script finished successful" ;;
          124|137) logger -t "MOS Plugins: $DISPLAY_NAME" "ERROR: start script timed out" ;;
          *)       logger -t "MOS Plugins: $DISPLAY_NAME" "ERROR: Running start script" ;;
        esac
        INSTALLATION_DONE=true
      fi
    fi
    if [ "$INSTALLATION_DONE" == "true" ] ; then
      logger -t "MOS Plugins: $DISPLAY_NAME" "Installation finished"
    fi
    unset INSTALLATION_DONE VERSION_PATH PACKAGE_PATH DISPLAY_NAME
  done
}

case "$1" in
  start)
    log_daemon_msg "Installing $DESC" "$NAME"
    install_plugins
    log_end_msg 0
  ;;
  stop)
    exit 0
  ;;
esac
